# conf/supervisor.conf.sample
[program:yourproject]
directory=/path/to/project/
command=uwsgi -x conf/uwsgi.xml
user=yourproject_user